# coreboot.service
[Unit]
Description=Conditionally Reboot the System

[Service]
Type=simple
ExecStart=/usr/local/sbin/coreboot
