# coreboot.service
[Unit]
Description=Conditional reboot service

[Service]
Type=simple
#Environment="DEBUG=1"
ExecStart=/usr/local/sbin/conreboot --daemon
